# Configs/Home/.config/autostart/bauh-tray.desktop
[Desktop Entry]
Exec=bauh-tray
Icon=
Name=bauh-tray
Path=
Terminal=False
Type=Application
