# Configs/Home/.config/autostart/pamac-tray-plasma.desktop
[Desktop Entry]
Exec=pamac-tray-plasma
Icon=
Name=pamac-tray-plasma
Path=
Terminal=False
Type=Application
